# linklandaemon.service
[Unit]
Description=Linklan Daemon service

[Service]
ExecStart=/usr/local/bin/linklandaemon

[Install]
WantedBy=multi-user.target
